# jarvis-telegram.service
[Unit]
Description=JARVIS Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mcde
WorkingDirectory=/home/mcde/whatsapp-claude
ExecStart=/usr/bin/node telegram.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jarvis-telegram

Environment=HOME=/home/mcde TZ=Europe/Berlin

# Sicherheits-Einstellungen
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/mcde/whatsapp-claude /home/mcde/.cache /tmp

[Install]
WantedBy=multi-user.target
